# lb2/api/api.service
[Unit]
Description=api

[Service]
ExecStart=/home/vagrant/main.py

[Install]
WantedBy=multi-user.target